# rclone/rclone-bisync.service
[Unit]
Description=rclone bisync for Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync /home/indervir/gdrive gdrive: \
  --transfers 8 \
  --checkers 16 \
  --conflict-resolve newer \
  --conflict-loser num \
  --conflict-suffix conflict \
  --recover \
  --resilient \
  --max-delete 50 \
  --check-sync true \
  --create-empty-src-dirs \
  --fix-case \
  --log-file /home/indervir/.local/share/rclone/bisync.log \
  --log-level INFO
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=default.target
